# donetick/rootfs/etc/services.d/donetick/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the donetick service with environment setup and process management
# ==============================================================================

# Set environment variables
# export VITE_APP_API_URL=http://host.docker.internal:2021

# Read configuration and export environment variables
export DT_TELEGRAM_TOKEN=$(bashio::config 'telegram_token')
export DT_PUSHOVER_TOKEN=$(bashio::config 'pushover_token')
export DT_DISABLE_SIGNUP=$(bashio::config 'disable_signup')


export DT_OAUTH2_CLIENT_ID=$(bashio::config 'oauth2_client_id')
export DT_OAUTH2_CLIENT_SECRET=$(bash::config 'oauth2_client_secret')
export DT_OAUTH2_REDIRECT_URI=$(bashio::config 'oauth2_redirect_uri')
export DT_OAUTH2_AUTH_URL=$(bashio::config 'oauth2_auth_url')
export DT_OAUTH2_TOKEN_URL=$(bashio::config 'oauth2_token_url')
export DT_OAUTH2_USER_INFO_URL=$(bashio::config 'oauth2_user_info_url')
export DT_OAUTH2_NAME=$(bashio::config 'oauth2_name')
export DT_JWT_SECRET=$(bashio::config 'jwt_secret')

cd /app
# Start donetick backend and save PID
./donetick-core &
PID1=$!

# Start nginx in the foreground as the second process
nginx -g 'daemon off;' &
PID2=$!

cleanup() {
    echo "Terminating processes..."
    kill $PID2
}

# Trap SIGINT & SIGTERM to clean up before exiting
trap cleanup SIGINT SIGTERM

# Wait for both processes to exit
wait $PID2

echo "Both processes have completed."